# root/etc/uci-defaults/luci-theme-lumos
#!/bin/sh

uci add luci themes
uci set luci.@themes[-1].name=LumOS
uci set luci.@themes[-1].mediaurlbase=/luci-static/lumos
uci commit luci

exit 0